# misc/systemd/dde-control-center.service.in
[Unit]
Description=DDE Control Center background service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-control-center -d
Type=dbus
BusName=org.deepin.dde.ControlCenter1
Slice=app.slice
